# deploy/vps/subscription-api.service
[Unit]
Description=Motion Video Subscription API (Deno)
After=network.target

[Service]
Type=simple
User=motionapi
WorkingDirectory=/opt/motion-video
EnvironmentFile=/opt/motion-video/deploy/vps/subscription-api.env
ExecStart=/usr/local/bin/deno run --allow-net --allow-env --env-file=/opt/motion-video/deploy/vps/subscription-api.env /opt/motion-video/supabase/functions/subscription-api/index.ts
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/motion-video

[Install]
WantedBy=multi-user.target
